# glpi-prod/docker/supervisord.conf
[supervisord]
nodaemon=true
[program:php-fpm]
command=docker-php-entrypoint php-fpm
[program:cron]
command=crond -f -L /var/log/cron.log
